# websocket.conf
[program:websocket-server]
command=php /path/to/your/websocket-server.php
directory=/path/to/your/project
autostart=true
autorestart=true
stderr_logfile=/var/log/websocket.err.log
stdout_logfile=/var/log/websocket.out.log
user=www-data 